# root/etc/services.d/radarr/run
#!/usr/bin/with-contenv bash

umask 022

cd /opt/Radarr || exit

exec \
	s6-setuidgid app mono --debug Radarr.exe \
	-nobrowser -data=/config
